# donspawn/mkosi.extra/mastodon-daily-task.timer
[Unit]
Description=Mastodon daily jobs

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
