# donspawn/mkosi.extra/set_up_mastodon.service
[Unit]
Description=set up mastodon on boot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set_up_mastodon.sh
