# torspawn/mkosi.extra/etc/systemd/system/tor.service.d/override.conf
[Service]
After=dnscrypt-proxy.service unbound.service
